# udev/10-android-tether.rules
# For Android tethering. 
ACTION=="add|remove", SUBSYSTEM=="net", RUN+="/usr/bin/systemctl restart systemd-networkd.service"
